# slash/usr/lib/systemd/system/pm2.service
[Unit]
Description=pm2 Service

[Service]
ExecStart=/bin/bash -c "/opt/pm2/extras/pm2-start"

[Install]
WantedBy=multiuser.target
