# slash/opt/pm2/extras/pm2.service
[Unit]
Description=pm2 Service

[Service]
ExecStart=/bin/bash -c "/opt/spark/extras/pm2-start"

[Install]
WantedBy=multiuser.target
